# system/etc/systemd/system/NetworkManager.service.d/override.conf
[Service]
ExecStop=/usr/bin/nmcli networking off
# Enforce a timeout to avoid hanging when shutting down
TimeoutStopSec=5
