# system/etc/systemd/system/bluetooth.service.d/override.conf
[Unit]
After=dbus.service

[Service]
ExecStop=/usr/bin/bluetoothctl power off
# Enforce a timeout to avoid hanging when shutting down
TimeoutStopSec=5
